# deploy/systemd/polymarket-bot.service
[Unit]
Description=Polymarket Trading Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/polymarket-bot
Environment=PYTHONPATH=/home/ubuntu/polymarket-bot
EnvironmentFile=/home/ubuntu/polymarket-bot/.env
ExecStart=/home/ubuntu/polymarket-bot/venv/bin/python -m src.main
Restart=always
RestartSec=10
TimeoutStopSec=30
# App logging can write to logs/bot.log via LOG_FILE; keep stdout in journald to avoid duplicates.
StandardError=append:/home/ubuntu/polymarket-bot/logs/bot-error.log

[Install]
WantedBy=multi-user.target
